# build/supervisor/laravel-horizon.conf
[program:horizon]
directory=/var/www/html
process_name=%(program_name)s
command=php artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=www-data
stdout_logfile=/var/www/html/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=120